# supervisord.conf
[supervisord]
nodaemon=true                 ; (start in foreground if true;default false)
user=root

[program:pgsync]
command=/usr/bin/pgsync listen
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
